# systemd/monitoring.service
[Unit]
Description=Monitoring script for test process

[Service]
ExecStart=/home/zhra/projects/monitoring/monitoring.sh
User=zhra
Restart=no
